# deploy/cloudflare/cloudflared.service.example
# cloudflared systemd service unit file template
#
# This file runs cloudflared as a persistent service that:
# - Automatically starts on reboot
# - Restarts on failure
# - Runs as an unprivileged user
#
# Installation:
#   sudo cp cloudflared.service.example /etc/systemd/system/cloudflared.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable cloudflared
#   sudo systemctl start cloudflared
#
# Commands:
#   sudo systemctl status cloudflared  # View status
#   sudo systemctl start cloudflared   # Start
#   sudo systemctl stop cloudflared    # Stop
#   sudo systemctl restart cloudflared # Restart
#   sudo journalctl -u cloudflared -f  # View logs

[Unit]
Description=Cloudflare Tunnel for RPi Drone Control
Documentation=https://developers.cloudflare.com/cloudflare-one/connections/connect-apps/
After=network.target
Wants=network-online.target

[Service]
Type=notify
User=pi
Group=pi
WorkingDirectory=/home/pi

# IMPORTANT: Customize these paths to match your setup
# Path to config file
Environment="CONFIG_PATH=/home/pi/.cloudflared/config.yml"

# Tunnel name (must match tunnel created with: cloudflared tunnel create)
Environment="TUNNEL_NAME=drone-control"

# Log level: debug, info, warn, error
Environment="LOG_LEVEL=info"

# The actual command to run
# Using tunnel name: (credentials auto-loaded from ~/.cloudflared/)
ExecStart=/usr/local/bin/cloudflared tunnel \
    --config %E/CONFIG_PATH \
    --loglevel %E/LOG_LEVEL \
    run %E/TUNNEL_NAME

# Auto-restart on failure
Restart=on-failure
RestartSec=10s

# Restart even if cloudflared exits cleanly (for high availability)
# Comment out if you prefer single-instance only
RestartForceExitStatus=0 1

# Process management
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudflared

# Security hardening
NoNewPrivileges=true
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/pi/.cloudflared

# Resource limits (adjust as needed)
# LimitNOFILE=65535
# LimitNPROC=512

# Timeout for service startup
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target

# After installation, verify with:
# sudo systemctl daemon-reload
# sudo systemctl enable cloudflared.service
# sudo systemctl start cloudflared.service
# sudo systemctl status cloudflared.service
# sudo journalctl -u cloudflared -f
