# AquaWing/deploy/systemd/rpi_high_level.service.example
# rpi_high_level systemd service unit file
#
# This file runs the RPi Drone Control backend as a persistent service.
#
# The backend FastAPI server hosts:
# - REST API at /api/*
# - WebSocket at /ws/*
# - Frontend dashboard at /
#
# Installation:
#   sudo cp rpi_high_level.service.example /etc/systemd/system/rpi_high_level.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable rpi_high_level
#   sudo systemctl start rpi_high_level
#
# Commands:
#   sudo systemctl status rpi_high_level
#   sudo systemctl start rpi_high_level
#   sudo systemctl stop rpi_high_level
#   sudo systemctl restart rpi_high_level
#   sudo journalctl -u rpi_high_level -f

[Unit]
Description=RPi High-Level Drone Control System
Documentation=https://github.com/your-repo/rpi_high_level
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/rpi_high_level

# Environment variables
Environment="PYTHONUNBUFFERED=1"
Environment="LOG_LEVEL=INFO"

# Method 1: Run main.py directly
ExecStart=/usr/bin/python3 /home/pi/rpi_high_level/main.py

# Method 2: Run as a module (alternative)
# ExecStart=/usr/bin/python3 -m backend.server

# Auto-restart on failure
Restart=on-failure
RestartSec=10s

# Process management
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rpi_high_level

# Security hardening
NoNewPrivileges=true
PrivateTmp=yes

# File descriptor limits
LimitNOFILE=65535

# Timeout for startup
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target

# After installation, verify with:
# sudo systemctl daemon-reload
# sudo systemctl enable rpi_high_level.service
# sudo systemctl start rpi_high_level.service
# sudo systemctl status rpi_high_level.service
# sudo journalctl -u rpi_high_level -f
